# saltstack/etc/cron.d/ocsinventory-agent
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/root/bin
49 7 * * * root /usr/local/bin/ocsinventory-agent --lazy > /dev/null 2>&1
