# saltstack/etc/cron.d/openvas-sync-plugins
# start plugin sync daily at midnight
0 0 * * * root /usr/sbin/openvas-nvt-sync-cron
